# deploy/systemd/komodo-marketmaker-gui.service
[Unit]
Description=Komodo Marketmaker GUI (Next.js with KCB layer)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/komodo-marketmaker-gui
EnvironmentFile=/opt/komodo-marketmaker-gui/.env.local
Environment=NODE_ENV=production
Environment=NEXT_TELEMETRY_DISABLED=1
ExecStart=/usr/bin/npm run start -- --hostname 0.0.0.0 --port 3000
Restart=on-failure
RestartSec=5
TimeoutStartSec=90
TimeoutStopSec=30
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal

# Hardening (lightweight, compatible with Next.js runtime)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/opt/komodo-marketmaker-gui /var/log

[Install]
WantedBy=multi-user.target
